# ubuntu/ptokax/PtokaX.conf
[supervisord]
nodaemon=true

[program:ptokax]
command=/opt/PtokaX/PtokaX -c /opt/PtokaX
