# ubuntu/ptokax/PtokaX.conf
[supervisord]
nodaemon=true

[program:PtokaX]
directory=/opt/PtokaX
command=/opt/PtokaX/PtokaX -c /opt/PtokaX
user=root
autostart=true
autorestart=true
