# ubuntu/postgres9.3/supervisord.conf
[supervisord]
nodaemon=true

[program:postgres]
user=postgres
command=/usr/lib/postgresql/9.3/bin/postgres -D /var/lib/postgresql/9.3/main --config_file=/etc/postgresql/9.3/main/postgresql.conf
autostart=true
autorestart=true
