# ubuntu/pgbouncer/pgbouncer.conf
[supervisord]
nodaemon=true

[program:pgbouncer]
#user=postgres
user=root
command=/usr/sbin/pgbouncer -u postgres /etc/pgbouncer/pgbconf.ini
autostart=true
autorestart=true
